# hdmi-switch/90-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/alsa/hdmi-switch.sh"
